# configs/wezterm/sysinfo-daemon/com.user.sysinfo-daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.sysinfo-daemon</string>

    <key>ProgramArguments</key>
    <array>
        <string>__HOME__/.local/bin/sysinfo-daemon</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardErrorPath</key>
    <string>/tmp/sysinfo-daemon.log</string>

    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
